# mothership-startup.service
[Unit]
Description=Update config

[Service]
Type=oneshot
User=root
ExecStart=/opt/microarch-custom/apply.sh

[Install]
WantedBy=multi-user.target
